# init/x86_64/macOS_10.12/com.shippable.reqKick.cron.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.shippable.reqKick.cron</string>

        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/bin/node /var/lib/shippable/{{UUID}}/reqKick/reqKick.app.js</string>
        </array>

        <key>Nice</key>
        <integer>1</integer>

        <key>RunAtLoad</key>
        <true/>

        <key>EnvironmentVariables</key>
        <dict>
            <key>STATUS_DIR</key>
            <string>{{STATUS_DIR}}</string>
            <key>SCRIPTS_DIR</key>
            <string>{{SCRIPTS_DIR}}</string>
            <key>REQEXEC_BIN_PATH</key>
            <string>{{REQEXEC_BIN_PATH}}</string>
            <key>RUN_MODE</key>
            <string>{{RUN_MODE}}</string>
        </dict>

        <key>StandardErrorPath</key>
        <string>{{LOG_DIR}}/reqKick_{{UUID}}.err</string>

        <key>StandardOutPath</key>
        <string>{{LOG_DIR}}/reqKick_{{UUID}}.out</string>
    </dict>
</plist>
